# dind-1/docker/supervisor/jenkins.conf
[program:jenkins]
command=/usr/local/bin/start_jenkins.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/jenkins.err.log
stdout_logfile=/var/log/jenkins.out.log